# assets/com.scott.codex-task-light.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.scott.codex-task-light.monitor</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>__CODEX_TASK_LIGHT_ROOT__/scripts/ensure_cx_running.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>/tmp/com.scott.codex-task-light.monitor.out</string>

    <key>StandardErrorPath</key>
    <string>/tmp/com.scott.codex-task-light.monitor.err</string>
  </dict>
</plist>
